# pslexplorersyn.crontab
# Sync crontab - put it into /etc/cron.d/

*/1 * * * * root cd /var/www/Pastel-explorer && /usr/bin/nodejs scripts/sync.js index update > /dev/null 2>&1
*/2 * * * * root cd /var/www/Pastel-explorer && /usr/bin/nodejs scripts/sync.js market > /dev/null 2>&1
*/5 * * * * root cd /var/www/Pastel-explorer && /usr/bin/nodejs scripts/peers.js > /dev/null 2>&1
